# home/.config/systemd/user/xdg-desktop-portal.service.d/kde-filechooser.conf
# KDE file picker needs KDE in XDG_CURRENT_DESKTOP when the portal starts.
[Service]
Environment="XDG_CURRENT_DESKTOP=Hyprland:KDE"
