# home/.config/cron/jobs.cron
# Managed by ~/dotfiles. Apply with: dotfiles cron

# Sync local Jira board markdown hourly during weekday working hours.
0 9-17 * * 1-5 /bin/zsh -lc 'jira sync' >/tmp/jira-sync.log 2>&1
